# RK3399Rootfs/AVLizardRootfs/etc/systemd/system/networking.service.d/reduce-timeout.conf
[Service]
TimeoutStartSec=30
TimeoutStopSec=5
